# server/sources/template-python/template.service
[Unit]
Description=Prospect news from {{SOURCE_NAME}} (Python)
ConditionPathExists=/home/sentinel/news/server
After=network.target

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/home/sentinel/news/server/sources/{{SOURCE_NAME}}
Environment=PYTHONPATH=/home/sentinel/news/server
ExecStart=/usr/local/bin/uv run main.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier={{SOURCE_NAME}}

# Python-specific settings
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target